# debian/DEBIAN/postinst
#!/bin/bash

systemctl --no-reload enable hobot-audio.service

exit 0

