# debian/DEBIAN/postrm
#!/bin/bash

systemctl disable hobot-audio.service
exit 0
